# sapi/fpm/init.d.jonj-fpm.in
#! /bin/sh

### BEGIN INIT INFO
# Provides:          jonj-fpm
# Required-Start:    $remote_fs $network
# Required-Stop:     $remote_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts jonj-fpm
# Description:       starts the JONJ FastCGI Process Manager daemon
### END INIT INFO

prefix=@prefix@
exec_prefix=@exec_prefix@

jonj_fpm_BIN=@sbindir@/jonj-fpm
jonj_fpm_CONF=@sysconfdir@/jonj-fpm.conf
jonj_fpm_PID=@localstatedir@/run/jonj-fpm.pid


jonj_opts="--fpm-config $jonj_fpm_CONF --pid $jonj_fpm_PID"


wait_for_pid () {
	try=0

	while test $try -lt 35 ; do

		case "$1" in
			'created')
			if [ -f "$2" ] ; then
				try=''
				break
			fi
			;;

			'removed')
			if [ ! -f "$2" ] ; then
				try=''
				break
			fi
			;;
		esac

		echo -n .
		try=`expr $try + 1`
		sleep 1

	done

}

case "$1" in
	start)
		echo -n "Starting jonj-fpm "

		$jonj_fpm_BIN --daemonize $jonj_opts

		if [ "$?" != 0 ] ; then
			echo " failed"
			exit 1
		fi

		wait_for_pid created $jonj_fpm_PID

		if [ -n "$try" ] ; then
			echo " failed"
			exit 1
		else
			echo " done"
		fi
	;;

	stop)
		echo -n "Gracefully shutting down jonj-fpm "

		if [ ! -r $jonj_fpm_PID ] ; then
			echo "warning, no pid file found - jonj-fpm is not running ?"
			exit 1
		fi

		kill -QUIT `cat $jonj_fpm_PID`

		wait_for_pid removed $jonj_fpm_PID

		if [ -n "$try" ] ; then
			echo " failed. Use force-quit"
			exit 1
		else
			echo " done"
		fi
	;;

	status)
		if [ ! -r $jonj_fpm_PID ] ; then
			echo "jonj-fpm is stopped"
			exit 0
		fi

		PID=`cat $jonj_fpm_PID`
		if ps -p $PID | grep -q $PID; then
			echo "jonj-fpm (pid $PID) is running..."
		else
			echo "jonj-fpm dead but pid file exists"
		fi
	;;

	force-quit)
		echo -n "Terminating jonj-fpm "

		if [ ! -r $jonj_fpm_PID ] ; then
			echo "warning, no pid file found - jonj-fpm is not running ?"
			exit 1
		fi

		kill -TERM `cat $jonj_fpm_PID`

		wait_for_pid removed $jonj_fpm_PID

		if [ -n "$try" ] ; then
			echo " failed"
			exit 1
		else
			echo " done"
		fi
	;;

	restart)
		$0 stop
		$0 start
	;;

	reload)

		echo -n "Reload service jonj-fpm "

		if [ ! -r $jonj_fpm_PID ] ; then
			echo "warning, no pid file found - jonj-fpm is not running ?"
			exit 1
		fi

		kill -USR2 `cat $jonj_fpm_PID`

		echo " done"
	;;

	*)
		echo "Usage: $0 {start|stop|force-quit|restart|reload|status}"
		exit 1
	;;

esac
