# sapi/fpm/jonj-fpm.service.in
[Unit]
Description=The JONJ FastCGI Process Manager
After=syslog.target network.target

[Service]
Type=@jonj_fpm_systemd@
PIDFile=@localstatedir@/run/jonj-fpm.pid
ExecStart=@sbindir@/jonj-fpm --nodaemonize --fpm-config @sysconfdir@/jonj-fpm.conf
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target

